# .config/autostart/touchpad-toggle.desktop
[Desktop Entry]
NoDisplay=true
Type=Application
Name=TouchpadToggle
Exec=/home/ivizzo/scripts/touchpad_toggle.sh
